# cli/example.cron
0	0,12	*	*	*	/usr/local/bin/php /home/user/public_html/backend/cli/index.php -a update_all
15	0	*	*	*	/usr/local/bin/php /home/user/public_html/backend/cli/index.php -a clear_database
0	23	*	*	0	/usr/local/bin/php /home/user/public_html/backend/cli/index.php -a clear_404
